# deploy/com.wc.refresh-all.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  本地定时刷三源(竞彩 sporttery + 欧盘 500 + Poly 7897代理)→ 推 HK 看板。
  取代旧的 com.wc.collect-zucai(它只推竞彩 raw)。
  ⚠️ 解释器必须含 httpx:本链经 tools/refresh_all.py → backend.sporttery(import httpx)强依赖 httpx,
     且 sporttery 文档化的 socks5 WAF 旁路还需 httpx[socks]——**禁用 /usr/bin/python3**(系统 3.9.6 无 httpx,
     启动即 ModuleNotFoundError、刷盘作业崩在 import)。这与已 stdlib 化、可用 /usr/bin/python3 的 backfill 不同。
  装之前: 把 PYTHON_WITH_HTTPX 改成装了 httpx 的 python3 绝对路径(如 /Users/heyining/bin/python3),
          把 REPO 改成本机仓库绝对路径, __看板密码__ 改成真实密码。
-->
<plist version="1.0"><dict>
  <key>Label</key><string>com.wc.refresh-all</string>
  <key>ProgramArguments</key>
  <array>
    <string>PYTHON_WITH_HTTPX</string>
    <string>REPO/tools/refresh_all.py</string>
    <string>--once</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>WC_INGEST_URL</key><string>http://18.166.71.60:8000</string>
    <key>WC_INGEST_PW</key><string>__看板密码__</string>
    <key>WC_HTTPS_PROXY</key><string>http://127.0.0.1:7897</string>
  </dict>
  <!-- 间隔取 1800(30min)与实机现状一致:刷三源对 HK/7897代理/Poly 的负载较低。
       design 原意是 300(5min);如需高频改回 300。模板对齐实机=照模板重装不会意外把频率 ×6。 -->
  <key>StartInterval</key><integer>1800</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardErrorPath</key><string>/tmp/wc-refresh-all.err</string>
  <key>StandardOutPath</key><string>/tmp/wc-refresh-all.out</string>
</dict></plist>
